# webvirt.supervisord.ini
[program:webvirtmgr]
command=/usr/share/webvirtmgr/bin/webvirt /usr/share/webvirtmgr/manage.py run_gunicorn -c /usr/share/webvirtmgr/conf/gunicorn.conf.py
directory=/usr/share/webvirtmgr
autostart=true
autorestart=true
logfile=/var/log/supervisor/webvirtmgr.log
log_stderr=true
user=apache

[program:webvirtmgr-console]
command=/usr/share/webvirtmgr/bin/webvirt /usr/share/webvirtmgr/console/webvirtmgr-console
directory=/usr/share/webvirtmgr
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/webvirtmgr-console.log
redirect_stderr=true
user=apache
